# deploy/supervisor/tornado_demo.ini
[program:tornado_demo]
environment=PYTHONPATH='/home/jiangtao/tornado-demo/.venv/lib/python3.6/site-packages/'
command=pipenv run python run.py
directory=/home/jiangtao/tornado-demo/app
user=nginx
stopwaitsecs=60
stopasgroup=true
killasgroup=true
autorestart=false
redirect_stderr=true
stdout_logfile=/var/log/jiangtao/tornado_demo.log
loglevel=info
